# deployment/systemd/pixel-workshop.service
[Unit]
Description=Pixel Workshop image tools
After=network-online.target
Wants=network-online.target

[Service]
Type=simple
User=pixelworkshop
Group=pixelworkshop
WorkingDirectory=/opt/pixel-workshop/current
Environment=NODE_ENV=production
Environment=PORT=3000
Environment=PUBLIC_BASE_PATH=/pictool/
ExecStart=/usr/bin/npm run start -- --hostname 127.0.0.1 --port 3000
Restart=on-failure
RestartSec=5
TimeoutStartSec=30
TimeoutStopSec=20
KillSignal=SIGTERM
NoNewPrivileges=true
PrivateTmp=true
ProtectSystem=full
ProtectHome=true

[Install]
WantedBy=multi-user.target
